# includes.chroot/etc/profile.d/qt-style-override.sh
export QT_STYLE_OVERRIDE=adwaita
export QT_QPA_PLATFORMTHEME=gtk3
